# mt500-playbook/files/etc/cron.d/mt500_watchdog
0 * * * * root python3 /root/watchdog.py && touch /root/ran
