# mt500-playbook/templates/etc/supervisor/conf.d/mt500.conf
[program:mt500]
directory=/home/{{ user }}/Desktop/pymt500/
command=python /home/{{ user }}/Desktop/pymt500/mt500.py
stderr_logfile=/var/log/mt500/error.log
stdout_logfile=/var/log/mt500/debug.log
autostart=true
autorestart=true
